# files/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock                   ; path to your socket file

[supervisord]
logfile=/var/log/supervisor/supervisord.log     ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=true                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                    ; default user
childlogdir=/var/log/supervisor/                ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock       ; use a unix:// URL  for a unix socket

[inet_http_server]
port=0.0.0.0:9001

[program:5000]
command=node app
directory=/app/docker/5000_port_hello_express
user=root
autostart=true
redirect_stderr=true
stderr_logfile = ./log/err.log
stdout_logfile = ./log/out.log

[program:5001]
command=node app
directory=/app/docker/5001_port_hello_express_css_view
user=root
autostart=true
redirect_stderr=true
stderr_logfile = ./log/err.log
stdout_logfile = ./log/out.log
